# n8n-ultimate-assistant/dev-bridge/dev-bridge.service
[Unit]
Description=Dev-Bridge — Onyx self-fix via Claude Code
After=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 %h/.dev-bridge/bridge.py
WorkingDirectory=%h/.dev-bridge
Restart=always
RestartSec=10

[Install]
WantedBy=default.target
