# n8n-ultimate-assistant/dev-bridge/vps-dev-bridge.service
[Unit]
Description=Dev-Bridge — Onyx self-fix on VPS (runs 24/7 with n8n)
After=network-online.target docker.service
Wants=docker.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/dev-bridge
EnvironmentFile=-/opt/dev-bridge/.env
ExecStart=/usr/bin/python3 /opt/dev-bridge/bridge.py
Restart=always
RestartSec=15
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
